# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid
user=root

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:caddy]
command=/usr/local/bin/start-caddy.sh
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uvicorn]
command=/opt/venv/bin/uvicorn backend.app:app --host %(ENV_BACKEND_HOST)s --port %(ENV_BACKEND_PORT)s --proxy-headers
directory=/app
autostart=true
autorestart=true
priority=20
# uvicorn drains the in-flight workflow on SIGTERM (see app.py `_shutdown`):
# worst case is run_script_timeout + wait_stopped_timeout +
# wait_running_timeout + mihomo_ready_timeout + 30s headroom (~960s with
# defaults). supervisord's default stopwaitsecs=10 would SIGKILL before
# the drain finishes and reintroduce the mid-`script_run` cancellation
# that the drain was added to prevent. Keep this in sync with
# TimeoutStopSec in deploy/mihomo-webui.container.
stopwaitsecs=1020
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app"
